# deploy/supervisor/whatsapp-bridge.conf
[program:nrsmarketing-whatsapp-bridge]
directory=/var/www/nrsmarketing/whatsapp-bridge
command=npm start
autostart=true
autorestart=true
startretries=3
user=www-data
environment=NODE_ENV="production",WHATSAPP_BRIDGE_TOKEN="change_me",WHATSAPP_BRIDGE_PORT="3000",WHATSAPP_SESSION_NAME="nrsmarketing"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nrsmarketing-whatsapp-bridge.log
